# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/opt/docker-postfix-s3/assets/postfix.sh

[program:rsyslog]
command=/usr/sbin/rsyslogd -n -c3